# websockets.conf
[program:websockets]
command=/usr/bin/php /artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=laravel-echo
